# scripts/systemd/poolmind-update.service
[Unit]
Description=PoolMind Auto-Update Service
After=network.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/PoolMind
ExecStart=/home/pi/PoolMind/scripts/deployment/auto-update.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
